# debian/contrail-webui-vendor/debian/contrail-webui-vendor.contrail-webui-jobserver.upstart
description "Contrail WebUI Job Server"
author "OpenContrail developers <dev@lists.opencontrail.org>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

chdir /var/lib/contrail-webui-vendor
respawn

script
  COMMAND="/usr/bin/node"
  CONF="/etc/contrail-webui-vendor/config.global.js"
  USER="contrail"
  OPTS="jobServerStart.js"

  # Allow override of command/conf and opts by /etc/default/daemon-name
  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi

  if ! [ -r "$CONF" ] ; then
    echo "Could not read ${CONF}: exiting"
    exit 0
  fi

  exec start-stop-daemon --start --chuid $USER:$USER --exec $COMMAND -- $OPTS

end script
