# rpm/contrail-openstack-config/rabbitmq-server.ini
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
redirect_stderr=true
stdout_logfile=/var/log/contrail/rabbit-server-supervisor-stdout.log
stderr_logfile=/var/log/contrail/rabbit-server-supervisor-stderr.log
priority=410
autostart=true
autorestart=true
startsecs=5
startretries=15
killasgroup=true
stopasgroup=true
stopsignal=TERM
exitcodes=0
