# debian/contrail/debian/database_files/kafka.ini
[program:kafka]
command=/usr/share/kafka/bin/kafka-server-start.sh /usr/share/kafka/config/server.properties
autostart=true                ; start at supervisord start (default: true)
killasgroup=false             ; SIGKILL the UNIX process group (def false)
