# rpm/contrail-web-core/contrail-webui.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization WebUI

sock_file=/tmp/supervisord_webui.sock
if [ -e "$sock_file" ]; then
    supervisorctl -s unix:///tmp/supervisord_webui.sock ${1} `basename ${0}`
else
    supervisorctl -s unix:///var/run/supervisord_webui.sock ${1} `basename ${0}`
fi

