# debian/ifmap-server/debian/ifmap-server.upstart
description "Contrail IFmap Server"
author "OpenContrail <dev@lists.opencontrail.org>"

start on runlevel [2345]
stop on runlevel [016]

chdir /var/run

script
  COMMAND=/usr/bin/ifmap-server
  USER=opencontrail

  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  exec start-stop-daemon --start --chuid $USER:$USER --exec $COMMAND --

end script
